# deploy/systemd/pia-run@.service
[Unit]
Description=Personal Investment Assistant — Monitor pipeline run (%i)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/ABS/PATH/TO/personal-investment-assistant
EnvironmentFile=/ABS/PATH/TO/personal-investment-assistant/.env
Environment=TZ=Europe/Rome
ExecStart=/ABS/PATH/TO/personal-investment-assistant/.venv/bin/pia-run --run-type %i
StandardOutput=append:/ABS/PATH/TO/personal-investment-assistant/logs/stdout.log
StandardError=append:/ABS/PATH/TO/personal-investment-assistant/logs/stderr.log
