# deploy/systemd/pia-bot.service
[Unit]
Description=Personal Investment Assistant — Advisor Telegram daemon
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/ABS/PATH/TO/personal-investment-assistant
EnvironmentFile=/ABS/PATH/TO/personal-investment-assistant/.env
Environment=TZ=Europe/Rome
ExecStart=/ABS/PATH/TO/personal-investment-assistant/.venv/bin/pia-bot
Restart=on-failure
StandardOutput=append:/ABS/PATH/TO/personal-investment-assistant/logs/bot-stdout.log
StandardError=append:/ABS/PATH/TO/personal-investment-assistant/logs/bot-stderr.log

[Install]
WantedBy=default.target
