# deploy/systemd/pia-pre-market.timer
[Unit]
Description=Personal Investment Assistant — pre-market run (08:00)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=pia-run@pre_market.service

[Install]
WantedBy=timers.target
