# deploy/launchd/com.personalinvestmentassistant.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personalinvestmentassistant.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/ABS/PATH/TO/personal-investment-assistant/.venv/bin/pia-bot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABS/PATH/TO/personal-investment-assistant</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/ABS/PATH/TO/personal-investment-assistant/logs/bot-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/ABS/PATH/TO/personal-investment-assistant/logs/bot-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>TZ</key>
        <string>Europe/Rome</string>
    </dict>
</dict>
</plist>
